# deploy/celery-printbuddy-cpu.service
[Unit]
Description=Celery CPU Worker — PrintBuddy (cpu_queue, concurrency=2)
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/printbuddy
Environment="DJANGO_SETTINGS_MODULE=printbuddy.settings.production"
ExecStart=/opt/printbuddy/venv/bin/celery \
    -A printbuddy worker \
    --queues=cpu_queue \
    --concurrency=2 \
    --loglevel=info \
    --logfile=/var/log/celery/printbuddy-cpu.log \
    --hostname=cpu@%%h
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
